# job.crontab
*/2 * * * * /Users/WorkMachine/GDrive/DTU/socialGraphs/twitter_bot/twitter_bot.py